# openrunlog.conf
description "open run log supervisord"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec bash /home/david/orl/start.bash >> /tmp/upstart.log 2>&1
